# worker/supervisord.conf
[program:celeryd]
command = python /home/dotcloud/current/manage.py celeryd --loglevel=DEBUG


